# files/etc/cron.daily/dehydrated
#!/usr/bin/env bash

# run dehydrated

/opt/dehydrated/dehydrated -c -g >/dev/null
