# files/etc/cron.daily/dehydrated
#!/usr/bin/env bash

# run dehydrated

{% for challenge_type in challenge_types %}/opt/dehydrated/dehydrated -f /etc/dehydrated/config_{{ challenge_type }} -c -g >/dev/null
{% endfor %}