# files/etc/cron.daily/dehydrated
#!/usr/bin/env bash

# run dehydrated

/opt/dehydrated/dehydrated -c >/dev/null
